# midscaled/install/midscaled.service
[Unit]
Description=Midscale WireGuard tunnel daemon
Documentation=https://github.com/midscale/midscaled
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/midscaled
Restart=always
RestartSec=5
StartLimitInterval=300
StartLimitBurst=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/midscaled /var/log/midscaled
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
DevicePolicy=closed
DeviceAllow=/dev/net/tun rw

# Environment file
EnvironmentFile=-/etc/midscaled/midscaled.env

# Logging
StandardOutput=journal
StandardError=journal

# Resource limits
LimitNOFILE=65536
LimitNPROC=128

[Install]
WantedBy=multi-user.target
